# inittab
::sysinit:/etc/init.d/rc
::askfirst:/bin/sh
::restart:/bin/init
::ctrlaltdel:/bin/reboot -f
::shutdown:/bin/umount -a -r
